# kasmvnc-docker-sysbox-runc/root/etc/s6-overlay/s6-rc.d/svc-de/run
#!/bin/bash

cd $HOME
/bin/bash /defaults/startwm.sh
